# certbot-img/crontab
# Run at 03:00 every day
0 3 * * * /opt/certbot_renew.sh
